# systemd/vservice.service
[Unit]
Description=vService Wrapper for vLLM
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/vllm/vllm.env
ExecStart=/opt/vllm/bin/vservice.sh
Restart=on-failure
RestartSec=5
User=vservice
Group=vservice
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
